# infra/systemd/zictia-api.service
[Unit]
Description=ZICTIA Customer Portal API
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=zictia
WorkingDirectory=/opt/zictia-portal/backend
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/opt/zictia-portal/backend/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zictia-api

[Install]
WantedBy=multi-user.target
